# deploy/frugl-api.service
[Unit]
Description=Frugl demo API (FastAPI + static frontend, same origin)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/home/claude/cp-advisor
# HOME is REQUIRED: extraction calls `claude -p`, which finds its OAuth creds under
# /home/claude. Without it every extraction silently 401s and falls back to degraded.
Environment=HOME=/home/claude
Environment=PATH=/usr/local/bin:/usr/bin:/bin
# The gate key + UTC expiry + leads path live here, NOT in the repo. See RUNBOOK.md.
EnvironmentFile=/etc/frugl.env
# `python3 -m` puts WorkingDirectory on sys.path so `dev_serve` imports, and StaticFiles
# (directory="frontend") resolves relative to it. Bind localhost only; the funnel fronts it.
ExecStart=/usr/bin/python3 -m uvicorn dev_serve:app --host 127.0.0.1 --port 47896
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
